# deployment/gunicorn.service
[Unit]
Description=Gunicorn daemon for Performance Review App
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/performance_review
EnvironmentFile=/home/ubuntu/performance_review/.env
ExecStart=/home/ubuntu/performance_review/venv/bin/gunicorn \
    --workers 2 \
    --bind unix:/home/ubuntu/performance_review/gunicorn.sock \
    --access-logfile /home/ubuntu/performance_review/logs/gunicorn-access.log \
    --error-logfile /home/ubuntu/performance_review/logs/gunicorn-error.log \
    performance_review.wsgi:application

[Install]
WantedBy=multi-user.target
